# etc/systemd/system/gre-tunnels%40.service
[Unit]
Description=GRE Tunnels for node %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/setup-gre.sh %i
Environment=GRE_TUNNELS_CONFIG=/etc/gre-tunnels/tunnels.conf
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gre-tunnels@%i

[Install]
WantedBy=multi-user.target
